# deploy/systemd/yuzu-gateway.service
[Unit]
Description=Yuzu Erlang/OTP Gateway
Documentation=https://github.com/Tr3kkR/Yuzu
After=network-online.target yuzu-server.service
Wants=network-online.target
# Restart-loop guard. The gateway is a release-built BEAM VM; if it
# crashes more than StartLimitBurst times within StartLimitIntervalSec
# (e.g. corrupt sys.config, missing certificate, port already bound),
# systemd marks the unit `failed` rather than spinning indefinitely.
# Keys live in [Unit] in systemd >= 230.
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=yuzu-gw
Group=yuzu-gw
ExecStart=/opt/yuzu_gw/bin/yuzu_gw foreground
Restart=on-failure
RestartSec=5

# File-descriptor headroom. The gateway aggregates many agent gRPC
# streams (one socket per agent under fanout) plus prometheus_httpd
# scrape sockets — the default 1024 limit caps useful fleet size at
# ~700 agents in practice. 65536 mirrors yuzu-server.service.
LimitNOFILE=65536

WorkingDirectory=/opt/yuzu_gw

# Erlang VM settings
Environment=ERL_CRASH_DUMP=/var/log/yuzu/gateway-crash.dump
Environment=RELX_REPLACE_OS_VARS=true

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/yuzu /opt/yuzu_gw
PrivateTmp=true

[Install]
WantedBy=multi-user.target
